# deployment/scraper/fyxed-scraper-huurwoningen.service
[Unit]
Description=Fyxed Wonen Scraper - Huurwoningen
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
EnvironmentFile=/etc/fyxedwonen-scraper.env
WorkingDirectory=/opt/Fyxedwonen
ExecStart=/bin/bash -lc 'SCRAPER_SOURCES=huurwoningen bash scraper/start.sh'
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7
ProtectSystem=full
ProtectHome=true
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target

